# test_data1/upstarttest.conf
start on lazy-start-40
script
	/test1.sh
end script
